# redis_exporter7001.service
[Unit]
Description=Redis Exporter for Redis node 7001
After=network.target, network-online.target

[Service]
User=nobody
ExecStart=/usr/local/bin/redis_exporter \
  --redis.addr redis://127.0.0.1:7001 \
  --web.listen-address :8001

Restart=always

[Install]
WantedBy=multi-user.target
