# redis_exporter7003.service
[Unit]
Description=Redis Exporter for Redis node 7003
After=network.target, network-online.target

[Service]
User=nobody
ExecStart=/usr/local/bin/redis_exporter \
  --redis.addr redis://127.0.0.1:7003 \
  --web.listen-address :8003

Restart=always

[Install]
WantedBy=multi-user.target
